# deploy/server.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  占位符 plist 模板 —— 别直接加载本文件。复制它、把下列占位符替换成你的真实值后，
  落到 ~/Library/LaunchAgents/<你的-label>.plist 再 bootstrap。

  占位符（下划线包裹、XML 安全，便于 sed 批量替换）：
    __LABEL__   LaunchAgent 标签，须与目标文件名一致，如 com.you.ccm-server
    __REPO__    仓库绝对路径，如 /Users/you/code/claude-chat-mobile
    __NODE__    node 绝对路径（command -v node），如 /opt/homebrew/bin/node
    __LOG__     日志绝对路径，如 /Users/you/Library/Logs/ccm-server.log

  一键替换（在仓库根执行，按需改 label）：
    sed -e 's#__LABEL__#com.you.ccm-server#' \
        -e 's#__REPO__#'"$PWD"'#' \
        -e 's#__NODE__#'"$(command -v node)"'#' \
        -e 's#__LOG__#'"$HOME/Library/Logs/ccm-server.log"'#' \
        deploy/server.plist.template > ~/Library/LaunchAgents/com.you.ccm-server.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <!-- 登录 shell 启动以保终端等价性：claude 的 PATH / 登录态与你终端一致 -->
    <string>cd __REPO__ &amp;&amp; exec __NODE__ server.js</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG__</string>
  <key>StandardErrorPath</key>
  <string>__LOG__</string>
</dict>
</plist>
