# deploy/tunnel.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  占位符 plist 模板 —— cloudflared 命名隧道常驻。别直接加载本文件。
  复制它、替换占位符后落到 ~/Library/LaunchAgents/<你的-tunnel-label>.plist 再 bootstrap。

  前置：已按 deployment.md「从零搭建 §1」建好命名隧道与 ~/.cloudflared/config.yml。
  这是标准 `cloudflared tunnel run` 配方；不含任何环境相关加固（如本机开全局代理时的
  edge 绕行），那类问题见 deployment.md 排错速查、按你自己的网络环境处理。

  占位符：
    __LABEL__         LaunchAgent 标签，如 com.you.ccm-tunnel
    __CLOUDFLARED__   cloudflared 绝对路径（command -v cloudflared），如 /opt/homebrew/bin/cloudflared
    __TUNNEL__        命名隧道名或 UUID（与 config.yml 的 tunnel 一致）
    __LOG__           日志绝对路径，如 /Users/you/Library/Logs/ccm-tunnel.log

  一键替换（在仓库根执行，按需改值）：
    sed -e 's#__LABEL__#com.you.ccm-tunnel#' \
        -e 's#__CLOUDFLARED__#'"$(command -v cloudflared)"'#' \
        -e 's#__TUNNEL__#your-tunnel-name#' \
        -e 's#__LOG__#'"$HOME/Library/Logs/ccm-tunnel.log"'#' \
        deploy/tunnel.plist.template > ~/Library/LaunchAgents/com.you.ccm-tunnel.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__CLOUDFLARED__</string>
    <string>tunnel</string>
    <string>run</string>
    <string>__TUNNEL__</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG__</string>
  <key>StandardErrorPath</key>
  <string>__LOG__</string>
</dict>
</plist>
